# usr/lib/udev/rules.d/51-nvidia-icd.rules
SUBSYSTEM=="pci", ACTION=="add", ENV{DRIVER}=="nvidia", RUN+="/usr/bin/change-nvidia-icd 'add'"
SUBSYSTEM=="pci", ACTION=="remove", ENV{DRIVER}=="nvidia", RUN+="/usr/bin/change-nvidia-icd 'remove'"
