# usr/lib/udev/rules.d/99-controller-wake.rules
ACTION=="add", SUBSYSTEM=="input", ENV{ID_INPUT_JOYSTICK}=="1", RUN+="/bin/systemctl start controller-wake@%E{DEVNAME}.service"
